# backend/cmdb-agent-go/launchd/com.cmdb.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cmdb.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cmdb-agent</string>
        <string>--config</string>
        <string>/etc/cmdb-agent/config.yaml</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/var/log/cmdb-agent/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/cmdb-agent/stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/var/lib/cmdb-agent</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
